# etc/s6-overlay/s6-rc.d/perkeepd-setup-identity/up
/etc/s6-overlay/scripts/perkeepd-setup-identity.sh
