# etc/s6-overlay/s6-rc.d/perkeepd-generate-config/up
/etc/s6-overlay/scripts/perkeepd-generate-config.sh
