# etc/s6-overlay/s6-rc.d/perkeepd/run
#!/command/execlineb -P
with-contenv
s6-setuidgid keeper
perkeepd -configfile /etc/perkeep/server-config.json
