# etc/s6-overlay/s6-rc.d/perkeep-identity-setup/up
/etc/s6-overlay/scripts/perkeep-identity-setup.sh
