# etc/s6-overlay/s6-rc.d/tailscaled/run
#!/command/execlineb -P
tailscaled --state=/var/lib/tailscale/tailscaled.state --socket=/var/run/tailscale/tailscaled.sock
